# com.ai-paper.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ai-paper.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/bytedance/ai-paper-tracker/papers-wrapper.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>12</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/bytedance/logs/paper-launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/bytedance/logs/paper-launchd.error.log</string>

    <!-- 即使错过了时间，唤醒后也执行 -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
